# scripts/scalp-bot-ctl.service
[Unit]
Description=Scalp Bot Telegram Controller
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /root/scalp-bot/scripts/telegram-ctl.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PM2_HOME=/root/.pm2
WorkingDirectory=/root/scalp-bot

[Install]
WantedBy=multi-user.target