# scripts/scalp-bot-ctl.service
[Unit]
Description=Scalp Bot Telegram Controller
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /root/tradingview-mcp-jackson/scripts/telegram-ctl.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
WorkingDirectory=/root/tradingview-mcp-jackson

[Install]
WantedBy=multi-user.target